# uploaddemo.service
[Unit]
Description=start videosite node

[Service]
ExecStart=/bin/bash /home/ubuntu/videosite/startup.sh

[Install]
WantedBy=multi-user.target